# service/cl.service
[Unit]
Description=none

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'pacman -Rns $(pacman -Qdtq) && pacman -Scc --noconfirm && fstrim --all'
